# supervisord.conf
[supervisord]
nodaemon=true

[program:openvpn]
command=openvpn --config server.conf
directory=/etc/openvpn
autorestart=true
